# docker/supervisord/websockets.conf
[program:websockets]
command=/usr/local/bin/php /var/www/apps/b2b/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=root
